# mysql/supervisord.conf
[program:mysqld]
command=/setup/mysql-start.sh
numprocs=1
autostart=true
autorestart=true
