# apache/supervisord.conf
[program:apache2]
command=/setup/apache-start.sh
numprocs=1
autostart=true
autorestart=true
